# roles/ecrc-motd/files/etc/update-motd.d/01-ecrc-header
#!/bin/sh

figlet -f slant ECRC : $(hostname -s)
